# service/version-api/backend/crontab
0 0  * * * root find /inference-api/tmp -type f -mtime +1 -exec rm {} +
# Don't remove the empty line at the end of this file. It is required to run the cron job